# 10-nook.rules
ACTION=="add", ATTRS{idVendor}=="2080", ATTRS{idProduct}=="0003", RUN+="/usr/local/bin/nook-script"
